# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
loglevel=info
# /tmp is always writable even with readOnlyRootFilesystem: true + emptyDir
pidfile=/tmp/supervisord.pid
childlogdir=/tmp
user=root

[program:otelcol]
command=/usr/local/bin/otelcol-contrib --config /etc/otelcol/config.yaml
autostart=true
autorestart=true
startretries=5
startsecs=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DEPLOYMENT_ENV="%(ENV_DEPLOYMENT_ENV)s",TEMPO_ENDPOINT="%(ENV_TEMPO_ENDPOINT)s",LOKI_ENDPOINT="%(ENV_LOKI_ENDPOINT)s",JAEGER_ENDPOINT="%(ENV_JAEGER_ENDPOINT)s"

[program:nginx]
command=/usr/local/bin/start-nginx.sh
autostart=true
autorestart=true
startretries=3
startsecs=2
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
